# trunk/allmon-client/conf/allmon-client.cron
###  scheduling-pattern [options] command [args]
###  the file is re-read every minute

###  Note!
###  Please, pay attention to how often your agents are being executed 
###  and how long they take for acquiring metrics from observed object!
###  If the time is longer than interval between calls, allmon will 
###  create more an more active agents threads, what can easily overload 
###  the mechanism and cause out of memory error and the scheduler process crash.

#sends heart-beat (ping message) to server - confirms that active agents scheduler and allmon client-side infrastructure works fine
* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HeartbeatAgent" "0"
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HeartbeatAgent" "20"
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HeartbeatAgent" "40"

#allmon agent self health check - use in debug mode only - avoid in production
* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.SelfHealthCheckAgent" "5"
* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.SelfHealthCheckAgent" "35"

#checks if google, wikipedia works
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HttpUrlCallAgent" shift setStrategy(getParamsString(0)); urlAddress = getParamsString(1); searchPhrase = getParamsString(2); contentType = getParamsString(3); urlParameters = getParamsString(4); checkingHost = getParamsString(5); checkName = getParamsString(6); instanceName = getParamsString(7); useProxy = Boolean.parseBoolean(getParamsString(8)) (true|false); requestMethod = getParamsString(9);
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HttpUrlCallAgent" "0" "org.allmon.client.agent.HttpUrlCallAgentBooleanStrategy" "http://www.wikipedia.org" "wikipedia" "text/html" "-" "www.wikipedia.com" "www.wikipedia.com-Checker" "original-wikipedia-instance" "true" "GET"
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HttpUrlCallAgent" "0" "org.allmon.client.agent.HttpUrlCallAgentBooleanStrategy" "http://www.google.com" "http" "text/html" "-" "www.google.com" "www.google.com-Checker" "original-google-instance" "true" "GET"

#checks if localhost active mq works
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.HttpUrlCallAgent" "0" "org.allmon.client.agent.HttpUrlCallAgentFullResponseTimeStrategy" "http://localhost:8161/admin/" "Broker" "text/html" "-" "localhost" "local-broker-check" "local-broker" "false" "POST"

#all local JVM - all JMX metrics (can be heavy if list of java processes is long)
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.JmxServerAgent" "0" ".*" ".*"
#checking local Active MQ JVM instance: all memory metrics
* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.JmxServerAgent" "2" ".*activemq.*" ".*java.lang:type=Memory.*"
#checking local Active MQ JVM instance: used memory and GC metrics
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.JmxServerAgent" "4" ".*activemq.*" "(.*java.lang:type=Memory.*used.*)|(.*java.lang:type=GarbageCollector.*)"
#checking local Active MQ JVM instance: JDK pools an OS memory allocations, GC stats, thread count, class loading stats
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.JmxServerAgent" "6" ".*activemq.*" "(.*java.lang:type=Memory.*used.*)|(.*java.lang:type=OperatingSystem.*)|(.*java.lang:type=GarbageCollector.*)|(.*java.lang:type=Threading.*ThreadCount.*)|(.*java.lang:type=ClassLoading.*)"
#checking all allmon processes: every two minutes -  JDK pools an OS memory allocations, GC stats, thread count, class loading stats
#*/2 * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.JmxServerAgent" "0" ".*allmon.*" "(.*java.lang:type=Memory.*used.*)|(.*java.lang:type=OperatingSystem.*)|(.*java.lang:type=GarbageCollector.*)|(.*java.lang:type=Threading.*ThreadCount.*)|(.*java.lang:type=ClassLoading.*)"

#SNMP for CPU
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.SnmpHostAgent" "6" "10.1.132.99" "-cpu"
#SNMP for process cpuTime and memory
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.SnmpHostAgent" "8" "10.1.132.99" "eclipse.exe"

#Shell call example: get task count (using: tasklist for windows)
#all processes
#*/10 * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.ShellCallAgent" "10" "tasklist /V /FO csv /NH | find /c \",\"" "[0-9]+"
#only allmon processes
#*/10 * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "org.allmon.client.agent.ShellCallAgent" "11" "tasklist /V /FO csv /NH | find /c \"allmon\"" "[0-9]+"
