# allmon-client/allmon-client.cron
#scheduling-pattern [options] command [args]
#the file is re-read every minute

#working simple call
#* * * * * C:\\Windows\\System32\\calc.exe

#does not work 
#* * * * * org.allmon.client.scheduler.RunableClass

#calling allmon metric collector
##* * * * * java org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.MetricCollector
##* * * * * java -cp bin;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.agent.MetricCollector
#* * * * * java -cp allmon-client.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.MetricCollector

#this is throwing an error - no class found
#* * * * * java -cp allmon-client.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.ShellCallAgentX
#check why it is not logging - because task.execute(); //XXX is not finishing!
#* * * * * java -cp allmon-client.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.ShellCallAgent

#checks if google works
#* * * * * java -cp allmon-client.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.URLCallAgent http://www.google.com \\d\\d\\d\\d

#check if iTropics
* * * * * java -cp allmon-client.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-camel-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/activemq-core-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-j2ee-management_1.0_spec-1.0.jar;../allmon-lib/lib/acivemq-5.2.0/geronimo-jms_1.1_spec-1.1.1.jar;../allmon-lib/lib/acivemq-5.2.0/optional/activemq-pool-5.2.0.jar;../allmon-lib/lib/acivemq-5.2.0/optional/commons-pool-1.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-jms-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/spring-tx-2.5.5.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xbean-spring-3.4.jar;../allmon-lib/lib/acivemq-5.2.0/optional/xmlpull-1.1.3.4d_b4_min.jar;../allmon-lib/lib/acivemq-5.2.0/commons-logging-1.1.jar;../allmon-lib/lib/spring-2.5.6.jar;../allmon-lib/lib/servlet-api-5.0.16.jar;../allmon-lib/lib/camel-1.6.0/jaxb-impl-2.1.6.jar;../allmon-lib/lib/camel-1.6.0/stax-api-1.0-2.jar;../allmon-lib/lib/camel-1.6.0/activation-1.1.jar;../allmon-lib/lib/camel-1.6.0/optional/commons-collections-3.2.1.jar;../allmon-lib/lib/camel-1.6.0/jaxb-api-2.1.jar;../allmon-lib/lib/camel-2.0-M2/camel-test-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-core-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jaxb-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jms-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-jpa-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-juel-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-2.0-M2.jar;../allmon-lib/lib/camel-2.0-M2/camel-spring-integration-2.0-M2.jar;../allmon-lib/lib/junit-4.4.jar;../allmon-lib/lib/log4j-1.2.15.jar;../allmon-lib/lib/cron4j-2.1.1.jar org.allmon.client.scheduler.AgentCallerMain org.allmon.client.agent.HttpUrlCallAgent http://lontd01/AdminConsole/statuscheck/default.aspx/CheckStatus \"Success\":true application/json;charset=utf-8 {'componentChecker':'TTC.iTropics.ComponentCheckers.ITropicsServiceComponentChecker,TTC.iTropics.ComponentCheckers'}

