# allmon-client/conf/allmon-client.cron
#scheduling-pattern [options] command [args]
#the file is re-read every minute

#checks if google, wikipedia works
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "0" "org.allmon.client.agent.HttpUrlCallAgent" setStrategy(getParamsString(0)); urlAddress = getParamsString(1); searchPhrase = getParamsString(2); contentType = getParamsString(3); urlParameters = getParamsString(4); checkingHost = getParamsString(5); checkName = getParamsString(6); instanceName = getParamsString(7); useProxy = Boolean.parseBoolean(getParamsString(8)) (true|false); requestMethod = getParamsString(9);
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "0" "org.allmon.client.agent.HttpUrlCallAgent" "org.allmon.client.agent.HttpUrlCallAgentBooleanStrategy" "http://www.wikipedia.org" "wikipedia" "text/html" "-" "www.wikipedia.com" "www.wikipedia.com-Checker" "original-wikipedia-instance" "true" "GET"
#* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "0" "org.allmon.client.agent.HttpUrlCallAgent" "org.allmon.client.agent.HttpUrlCallAgentBooleanStrategy" "http://www.google.com" "http" "text/html" "-" "www.google.com" "www.google.com-Checker" "original-google-instance" "true" "GET"

#all local JVM - all JMX metrics
* * * * * java:org.allmon.client.scheduler.AgentCallerMain#main "0" "org.allmon.client.agent.JmxServerAgent" "" ""
